# supervisor/app.conf
[program:node]
priority=5
directory=/app/code
environment=HOME=/app/code
command=node /app/code/api/index.js
user=cloudron
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mongo]
priority=5
directory=/app/code
environment=HOME=/app/code
command=/app/code/mongo/bin/mongod --dbpath /app/data/db --config /app/data/mongod.conf
user=cloudron
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


